# systemd/wh-reconcile.service
[Unit]
Description=Reconcile Warehouse services to their latest release on boot
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
# Replace <user> with the non-root account that owns /opt/wh and is in the
# docker group (the same user the runners run as).
User=<user>
ExecStart=/opt/wh/wh-autodeploys/scripts/reconcile.sh

[Install]
WantedBy=multi-user.target
